# scripts/pi-detect.service
[Unit]
Description=Pi-Detect Object Detection Livestream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-detect
Environment=PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/pi
ExecStartPre=/bin/bash -c "sudo mount --bind /home/pi/tmp /tmp || true"
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /home/pi/pi-detect/main.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/pi-detect/logs/service.log
StandardError=append:/home/pi/pi-detect/logs/service.log

[Install]
WantedBy=multi-user.target
